# trees/99-patches/etc/rc.d/rc.local
#!/bin/bash
#
# /etc/rc.d/rc.local:  Local system initialization script.
#
# Put any local startup commands in here.  Also, if you have
# anything that needs to be run at shutdown time you can
# make an /etc/rc.d/rc.local_shutdown script and put those
# commands in there.

# Run some deferred setup tasks if necessary
if [ -x /etc/rc.d/rc.firstrun ]; then
    /etc/rc.d/rc.firstrun
    chmod a-x /etc/rc.d/rc.firstrun
fi
# Set a random password and create new user
/bin/autopw startup
